# deploy/systemd/mem0-api.service
[Unit]
Description=mem0 API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/home/appuser/project
EnvironmentFile=/home/appuser/project/.env
Environment=PYTHONPATH=/home/appuser/project/.deps
ExecStart=/usr/bin/python3 -m uvicorn main:app --host ${MEM0_HOST} --port ${MEM0_PORT}
Restart=always
RestartSec=3
TimeoutStopSec=15
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
