# deploy/systemd/mem0-api.service
[Unit]
Description=mem0 API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vhserver
Group=vhserver
WorkingDirectory=/home/vhserver/bt
EnvironmentFile=/home/vhserver/bt/.env
ExecStart=/home/vhserver/bt/.venv/bin/python -m uvicorn main:app --host ${MEM0_HOST} --port ${MEM0_PORT}
Restart=always
RestartSec=3
TimeoutStopSec=15
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
